# assets/systemd/ming-transaction-health.service
[Unit]
Description=Ming OS transactional update health confirmation
After=local-fs.target dbus.service NetworkManager.service ming-transaction-reconcile.service
Before=display-manager.service graphical.target
OnFailure=ming-transaction-rollback-reboot.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ming-transaction-health confirm-active
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
