# assets/systemd/display-manager.service.d/20-ming-transaction-health.conf
[Unit]
Requires=ming-transaction-health.service
After=ming-transaction-reconcile.service
After=ming-transaction-health.service
